# config/includes.chroot/etc/skel/.config/autostart/boomerbox-firstboot.desktop
[Desktop Entry]
Type=Application
Name=BoomerBox First Boot Setup
Exec=/usr/local/bin/boomerbox-firstboot
X-GNOME-Autostart-enabled=true
NoDisplay=true
